# launchd/com.openclaw.session-janitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.session-janitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/boris/.openclaw/workspace/tools/session-janitor/session-janitor.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/boris</string>
  </dict>

  <key>Nice</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/boris/.openclaw/workspace/tools/session-janitor/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/boris/.openclaw/workspace/tools/session-janitor/logs/launchd-stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
